# rules/sigma/linux/file_event_lnx_cron_persistence.yml
title: Cron Job Persistence
id: 4e5f6071-8293-44a5-8eb6-3f4a5b6c7d84
status: experimental
description: >
  Detects creation of cron entries in system cron directories or user crontab
  spools, a classic scheduled-task persistence mechanism. Package installs also
  drop cron files, so intended for the Advanced tier with tuning.
references:
  - https://attack.mitre.org/techniques/T1053/003/
author: rustinel-rules
date: 2026-06-04
tags:
  - attack.persistence
  - attack.execution
  - attack.t1053.003
logsource:
  category: file_event
  product: linux
detection:
  selection_dir:
    TargetFilename|contains:
      - '/etc/cron.d/'
      - '/etc/cron.daily/'
      - '/etc/cron.hourly/'
      - '/etc/cron.weekly/'
      - '/etc/cron.monthly/'
      - '/var/spool/cron/'
  selection_file:
    TargetFilename|endswith: '/etc/crontab'
  condition: selection_dir or selection_file
falsepositives:
  - Package installation and legitimate scheduled administrative tasks.
level: medium
rustinel:
  telemetry:
    - file_event
  expected_false_positive_level: medium
  test_status: none
