# rules/sigma/linux/file_event_lnx_systemd_persistence.yml
title: Systemd Unit Persistence
id: 3d4e5f60-7182-4394-9da5-2e3f4a5b6c73
status: experimental
description: >
  Detects creation of a systemd service, timer or socket unit in a system or
  per-user unit directory, a durable persistence mechanism. Noisy on hosts where
  packages are installed/updated, so intended for the Advanced tier with tuning.
references:
  - https://attack.mitre.org/techniques/T1543/002/
author: rustinel-rules
date: 2026-06-04
tags:
  - attack.persistence
  - attack.privilege-escalation
  - attack.t1543.002
logsource:
  category: file_event
  product: linux
detection:
  selection_path:
    TargetFilename|contains:
      - '/etc/systemd/system/'
      - '/lib/systemd/system/'
      - '/usr/lib/systemd/system/'
      - '/.config/systemd/user/'
  selection_ext:
    TargetFilename|endswith:
      - '.service'
      - '.timer'
      - '.socket'
  condition: selection_path and selection_ext
falsepositives:
  - Package installation and updates that ship systemd units.
level: medium
rustinel:
  telemetry:
    - file_event
  expected_false_positive_level: medium
  test_status: atomic
