# packs/windows/essential/pack.yml
name: Windows Essential
id: windows-essential
description: Low-noise, high-confidence Windows detections for Rustinel. Safe default pack.
os: windows
level: essential
pack_schema_version: 1
requires_rustinel: ">=1.0.2"
default: true
expected_false_positive_level: low
status: experimental
license: DRL-1.1
extends: []
attack_coverage:
  - T1059.001
  - T1027
  - T1003.001
  - T1490
  - T1218.010
  - T1218.005
  - T1566.001
  - T1059
  - T1070.001
  - T1562.001
  - T1003.002
  - T1003.003
  - T1548.002
  - T1053.005
  - T1105
telemetry_requirements:
  - process_creation
  - registry_event
  - task_creation
  - file_scan
test_status: mixed
sources:
    manual:
      - https://attack.mitre.org/
      - https://lolbas-project.github.io/
rules:
    has:
        sigma:
          - 7f3a1c2e-4b5d-4e6f-8a90-1b2c3d4e5f60   # Suspicious Encoded PowerShell Command Line (sigma)   # Suspicious Encoded PowerShell Command Line
          - 7d2a8f31-4b6c-49e0-a1f8-3c5d9e0b2a05   # LSASS Memory Dump via comsvcs.dll MiniDump (sigma)
          - 1e9b4d68-2c7a-4f93-8b15-6d0a2e4c1b04   # Volume Shadow Copy Deletion (sigma)
          - 9a0c3e75-6d8b-4f12-bc34-1e7a5d9c0b06   # Regsvr32 Remote Scriptlet Execution / Squiblydoo (sigma)
          - 3c6f9b22-7e1a-4d85-9c0b-2a8d4e6f1b07   # Mshta Remote or Inline Script Execution (sigma)
          - c4a7f2e9-5d3b-4810-a6c1-9e0f3b7d2a03   # Office Application Spawning a Command Shell (sigma)
          - 5e8b1d44-9c2a-4f67-8d03-7b1a6e3c9d08   # Windows Event Log Cleared via Command Line (sigma)
          - b2c3d4e5-8f90-4123-8b4c-5d6e7f801a11   # Microsoft Defender Tampering via Registry (sigma)
          - e7a1b9c2-4d6f-4a83-b2e5-1c0d9f3a4b20   # Registry Hive Dump via reg.exe save (sigma)
          - f2b8c4d1-6e90-4a27-9d3b-5a1c2e7f8b21   # UAC Bypass via Auto-Elevating LOLBin (sigma)
          - a3c7e9b4-2f81-4d56-8c0a-6b2d1f4e9c22   # Active Directory Database (NTDS.dit) Extraction (sigma)
          - b4d8f0c5-3a92-4e67-9d1b-7c3e2f5a0d23   # WDigest Cleartext Credential Caching Enabled (sigma)
          - 3a1cdc64-a2a1-444b-a04b-32a3f1d09a99   # Scheduled Task Created With Suspicious Action (sigma)
        yara:
          - yara-win-mimikatz-strings              # Mimikatz credential-dumping strings (yara)
        ioc:
          - ioc-eicar-test                         # EICAR safe end-to-end IOC test set (ioc)                         # EICAR safe end-to-end IOC test set (ioc)                         # EICAR safe end-to-end IOC test set (ioc)
