# rules/sigma/windows/task_creation_win_susp_scheduled_task_action.yml
title: Scheduled Task Created With Suspicious Action
id: 3a1cdc64-a2a1-444b-a04b-32a3f1d09a99
status: experimental
description: >
  Detects newly created scheduled tasks whose action launches a script host,
  command interpreter or LOLBin together with suspicious payload indicators such
  as encoded PowerShell, inline script protocols, remote URLs or user-writable
  paths. This is a tight low-noise variant for the Essential tier.
references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://lolbas-project.github.io/
author: rustinel-rules
date: 2026-06-04
tags:
  - attack.execution
  - attack.persistence
  - attack.privilege-escalation
  - attack.t1053.005
  - attack.t1059
  - attack.t1105
logsource:
  category: task_creation
  product: windows
detection:
  selection_action:
    TaskContent|contains:
      - 'powershell.exe'
      - 'pwsh.exe'
      - 'cmd.exe'
      - 'wscript.exe'
      - 'cscript.exe'
      - 'mshta.exe'
      - 'rundll32.exe'
      - 'regsvr32.exe'
      - 'bitsadmin.exe'
      - 'certutil.exe'
  selection_payload:
    TaskContent|contains:
      - ' -enc '
      - ' -encodedcommand '
      - ' /enc '
      - 'http://'
      - 'https://'
      - 'javascript:'
      - 'vbscript:'
      - '\Users\'
      - '\Temp\'
      - '\AppData\'
      - '%TEMP%'
      - '%APPDATA%'
  condition: selection_action and selection_payload
falsepositives:
  - Rare administrative tasks that intentionally run scripts from user-writable locations or fetch remote content.
level: high
rustinel:
  telemetry:
    - task_creation
  expected_false_positive_level: low
  test_status: manual
  test_reason: Task creation telemetry did not fire reliably on hosted Windows runners with engine v1.1.1.
